# lego-scraper/doc/crontab
#!/usr/bin/env bash
45 3 * * * php -S lego-scrape.php >/dev/null 2>&1
